# deploy/imac-backup.timer
[Unit]
Description=Run daily backups for the IMAC mongodb

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
